# systemd/bugbounty-agent.service
[Unit]
Description=Liminal - Autonomous Security Scanner
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=liminal
WorkingDirectory=/opt/liminal
EnvironmentFile=/opt/liminal/.env
ExecStart=/opt/liminal/venv/bin/liminal scan \
    --config /opt/liminal/config/targets.yaml \
    --log-file /var/log/liminal/scan.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liminal
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
